# systemd/user/flameshot.service
[Unit]
Description=Screenshot daemon

[Service]
ExecStart=/usr/bin/flameshot

[Install]
WantedBy=default.target
